# personio-server.service.example
[Unit]
Description=personio-server
After=network.target

[Service]
ExecStart=/usr/local/bin/personio-server
Restart=always

[Install]
WantedBy=graphical-session.target